# supervisor.conf
[supervisord]
nodaemon=true

[program:styloko]
command="/scripts/start-styloko"

[program:org]
command="/scripts/start-org"
